# udev-rules/plugable_sound.rules
ACTION=="add", ATTRS{idVendor}=="17e9", ATTRS{idProduct}=="6011", RUN+="/usr/local/bin/activate_plugable_sound"
