# udev-rules/webcam_C920.rules
ACTION=="add", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="082d", RUN+="/usr/local/bin/activate_webcam_C920"
